# @fsystem/src/etc.skel/.config/autostart/RSp - rsp.wireguardConnect.desktop
[Desktop Entry]
Exec=/usr/local/bin/rsp.wireguardConnect
Name=RSp - rsp.wireguardConnect
Type=Application
Version=1.0
X-RSp-Need-Tray=true
