# rsp-fsystem/src/etc.skel/.config/autostart/RSp - rsp.disableIPV6.desktop
[Desktop Entry]
Exec=rsp.disableIPV6
Name=RSp - rsp.disableIPV6
Type=Application
Version=1.0
